# sources/root/etc/udev/rules.d/backlight.rules
ACTION=="add", SUBSYSTEM=="backlight", RUN+="/usr/bin/chgrp video $sys$devpath/brightness", RUN+="/usr/bin/chmod g+w $sys$devpath/brightness"
