# sources/files/home/.config/systemd/user.control/mpd.service.d/data-dir.conf
[Service]
ExecStartPre=/usr/bin/mkdir -p %h/.local/share/mpd
