# sources/root/etc/profile.d/bash-xdg.sh
. ~/.config/bash/bash_profile
. ~/.config/bash/bashrc
